# etc/service/minidlna/run
#!/bin/sh

. /etc/container_environment.sh

exec setuser 1000 /usr/sbin/minidlnad -v -f /ertc/minidlna.conf
