# etc/service/minidlna/run
#!/bin/sh

. /etc/container_environment.sh

exec /usr/sbin/minidlnad -d -v -f /etc/minidlna.conf
